# sparse/lib/udev/rules.d/999-bluetooth-rfkill.rules
SUBSYSTEM=="rfkill", ATTR{type}=="bluetooth", RUN+="/usr/bin/droid/bluetooth-rfkill.sh"
